# app/default/savedsearches.conf
############################################################################
# Pre-wired saved searches for the splunk-sigma app.
#
# Each search runs a slice of the bundled rules against a common Windows
# telemetry source. Enable/schedule them in Splunk Web as needed.
############################################################################

[sigma_all_process_creation]
description = Run every bundled Sigma rule against Sysmon process-creation events (EventCode=1).
search = index=* sourcetype=XmlWinEventLog EventCode=1 | sigma rules="*"
dispatch.earliest_time = -15m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 0

[sigma_credential_access]
description = Detect credential-access techniques (LSASS dumps, mimikatz).
search = index=* sourcetype=XmlWinEventLog EventCode=1 | sigma rules="attack:T1003.001"
dispatch.earliest_time = -15m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 0
alert.severity = 5
alert.track = 1

[sigma_persistence]
description = Detect persistence techniques (scheduled tasks, Run keys).
search = index=* sourcetype=XmlWinEventLog EventCode=1 | sigma rules="attack:T1053.005 OR attack:T1547.001"
dispatch.earliest_time = -15m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 0

[sigma_lateral_rdp_external]
description = Detect RDP logons from non-RFC1918 source IPs.
search = index=* sourcetype=WinEventLog:Security EventCode=4624 | sigma rules="attack:T1021.001"
dispatch.earliest_time = -15m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 0
alert.severity = 4

[sigma_lolbin_downloads]
description = Detect living-off-the-land downloads (certutil, etc).
search = index=* sourcetype=XmlWinEventLog EventCode=1 | sigma rules="attack:T1105"
dispatch.earliest_time = -15m
dispatch.latest_time = now
cron_schedule = */5 * * * *
enableSched = 0
